# weather.timer
[Unit]
Description=Run weather hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=multi-user.target
